# source/usr/share/doc/bung/examples/ansible/bung/files_for_systemd/etc/systemd/system/bung_all.slice
[Unit]
Description=bung super_bu.sh all job slice
DefaultDependencies=no
Before=slices.target

[Slice]
CPUQuota=25%
