# source/usr/share/doc/bung/examples/ansible/bung/files_for_systemd/etc/systemd/system/bung_all.service
[Unit]
Description=bung super_bu.sh all job

[Service]
Type=simple
ExecStart=/usr/bin/super_bu -c all
User=root
Group=root
Slice=bung_all.slice
